# blue-team-study-portal/data/json/pages/romcom-htb.json
{
  "id": "forensics labs/romcom (htb)",
  "title": "RomCom (HTB)",
  "slug": "romcom-htb",
  "relPath": "Forensics Labs/RomCom (HTB).md",
  "level": "L2",
  "curriculumLevelId": "l2",
  "curriculumSectionId": "l2-dfir",
  "curriculumSectionTitle": "Digital Forensics (DFIR)",
  "tags": [
    "Windows",
    "Network",
    "DFIR",
    "Threat Intel",
    "Persistence"
  ],
  "mitre": [
    "T1547.009"
  ],
  "parentId": "forensics labs",
  "children": [],
  "relatedLinks": [],
  "backlinks": [
    "forensics labs"
  ],
  "breadcrumbs": [
    {
      "id": "forensics labs",
      "title": "Forensics Labs",
      "slug": "forensics-labs"
    },
    {
      "id": "forensics labs/romcom (htb)",
      "title": "RomCom (HTB)",
      "slug": "romcom-htb"
    }
  ],
  "headings": [
    {
      "level": 1,
      "text": "RomCom (HTB)"
    },
    {
      "level": 2,
      "text": "Scenario"
    },
    {
      "level": 3,
      "text": "1- What is the CVE assigned to the WinRAR vulnerability exploited by the RomCom threat group in 2025?"
    },
    {
      "level": 3,
      "text": "CVE-2025-8088"
    },
    {
      "level": 3,
      "text": "2- What is the nature of this vulnerability?"
    },
    {
      "level": 3,
      "text": "Path Traversal"
    },
    {
      "level": 3,
      "text": "3- What is the name of the archive file under Susan's documents folder that exploits the vulnerability upon opening the archive file?"
    },
    {
      "level": 3,
      "text": "Pathology-Department-Research-Records.rar"
    },
    {
      "level": 3,
      "text": "4- When was the archive file created on the disk?"
    },
    {
      "level": 3,
      "text": "2025-09-02 08:13:50"
    },
    {
      "level": 3,
      "text": "5- When was the archive file opened?"
    },
    {
      "level": 3,
      "text": "2025-09-02 08:14:04"
    },
    {
      "level": 3,
      "text": "6- What is the name of the decoy document extracted from the archive file, meant to appear legitimate and distract the user?"
    },
    {
      "level": 3,
      "text": "GenotypingResultsB57Positive.pdf"
    },
    {
      "level": 3,
      "text": "7- What is the name and path of the actual backdoor executable dropped by the archive file?"
    },
    {
      "level": 3,
      "text": "C:\\Users\\Susan\\Appdata\\Local\\ApbxHelper.exe"
    },
    {
      "level": 3,
      "text": "8- The exploit also drops a file to facilitate the persistence and execution of the backdoor. What is the path and name of this file?"
    },
    {
      "level": 3,
      "text": "C:\\Users\\Susan\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\Display Settings.lnk"
    },
    {
      "level": 3,
      "text": "9- What is the associated MITRE Technique ID discussed in the previous question?"
    },
    {
      "level": 3,
      "text": "T1547.009"
    },
    {
      "level": 3,
      "text": "10- When was the decoy document opened by the end user, thinking it to be a legitimate document?"
    },
    {
      "level": 3,
      "text": "2025-09-02 08:15:05"
    }
  ],
  "overview": "Susan works at the Research Lab in Forela International Hospital. A Microsoft Defender alert was received from her computer, and she also mentioned that while extracting a document from the received file, she received tons of errors, but the document opened just fine. According to the latest threat intel feeds, WinRAR is being exploited in the wild to gain i",
  "study": {
    "overview": "Susan works at the Research Lab in Forela International Hospital. A Microsoft Defender alert was received from her computer, and she also mentioned that while extracting a document from the received file, she received tons of errors, but the document opened just fine. According to the latest threat intel feeds, WinRAR is being exploited in the wild to gain i",
    "whySoc": "Use this note as an analyst reference: understand the concept, identify the relevant telemetry, pivot through evidence, and document the decision path clearly.",
    "keyConcepts": [
      "RomCom (HTB)",
      "Scenario",
      "1- What is the CVE assigned to the WinRAR vulnerability exploited by the RomCom threat group in 2025?",
      "CVE-2025-8088",
      "2- What is the nature of this vulnerability?",
      "Path Traversal",
      "3- What is the name of the archive file under Susan's documents folder that exploits the vulnerability upon opening the archive file?",
      "Pathology-Department-Research-Records.rar",
      "4- When was the archive file created on the disk?",
      "2025-09-02 08:13:50",
      "5- When was the archive file opened?",
      "2025-09-02 08:14:04"
    ],
    "workflow": [
      "Identify the user, host, process, network flow, or identity object involved.",
      "Collect the highest-value logs first and preserve timestamps and assumptions.",
      "Pivot through related notes, telemetry sources, and MITRE references.",
      "Document findings, gaps, containment options, and detection improvements."
    ],
    "splunkExamples": [],
    "importantCommands": [],
    "tools": [],
    "tips": []
  }
}